# services/bootlog.service
[Unit]
Description=Boot diagnostics and self-healing
After=local-fs.target boot-firmware.mount network.target
Wants=local-fs.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/bootlog.sh
RemainAfterExit=yes
TimeoutStartSec=60

[Install]
WantedBy=multi-user.target
