# services/buffer-cleaner.service
[Unit]
Description=Buffer and memory cleaner

[Service]
Type=oneshot
ExecStart=/usr/local/bin/buffer-cleaner.sh
