# services/buffer-cleaner.timer
[Unit]
Description=Run buffer cleaner every 5 minutes

[Timer]
OnBootSec=120
OnUnitActiveSec=300

[Install]
WantedBy=timers.target
