# services/oxigotchi-splash.service
[Unit]
Description=Oxigotchi Boot Splash
DefaultDependencies=no
Before=pwnagotchi.service
After=local-fs.target sysinit.target

[Service]
Type=oneshot
ExecStart=/home/pi/.pwn/bin/python3 /usr/local/bin/oxigotchi-splash.py /etc/pwnagotchi/custom-plugins/faces/awake.png
ExecStop=/home/pi/.pwn/bin/python3 /usr/local/bin/oxigotchi-splash.py /etc/pwnagotchi/custom-plugins/faces/shutdown.png
RemainAfterExit=yes
TimeoutStartSec=30
TimeoutStopSec=30

[Install]
WantedBy=sysinit.target
