# services/pisugar-watchdog.service
[Unit]
Description=PiSugar I2C watchdog — reconnects pisugar-server when MCU wakes
After=pisugar-server.service

[Service]
Type=oneshot
ExecStart=/usr/local/bin/pisugar-watchdog.sh
