# services/usb0-fallback.service
[Unit]
Description=USB0 IP fallback
After=NetworkManager.service

[Service]
Type=oneshot
ExecStart=/usr/local/bin/usb0-fallback.sh
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
